# ssh/supervisord.conf
[supervisord]
nodaemon=true

[program:gogs]
command=/etc/init.d/ssh start 
