# memcache/supervisord.conf
[supervisord]
nodaemon=true

[program:memcache]
command=memcache -d
