# gogs/supervisord.conf
[supervisord]
nodaemon=true

[program:gogs]
command=/opt/gogs web 
autostart=true
